# scripts/daedalus.service
[Unit]
Description=Daedalus Ground Control Station v1.0 (Team #1043)
Documentation=https://github.com/Crazyexs/DDL-JS
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/DDL-JS

# Use the venv Python directly so no activation script is needed
ExecStart=/home/pi/DDL-JS/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8080 --log-level info

# Restart automatically if it crashes, wait 5 seconds between restarts
Restart=always
RestartSec=5

# Log everything to the system journal (view with: journalctl -u daedalus -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daedalus

# Give the process time to start before marking it as failed
TimeoutStartSec=30

# Environment variables (no TTY = headless auto-port detection kicks in)
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
